# infra/systemd/pj-general-sync.timer
[Unit]
Description=Run pj-general source sync every six hours

[Timer]
OnCalendar=*-*-* 00,06,12,18:00:00
Persistent=true
RandomizedDelaySec=5min
Unit=pj-general-sync.service

[Install]
WantedBy=timers.target
